# deploy/tencent-cloud/world-cup-prediction-daily.timer
[Unit]
Description=Run World Cup Prediction Daily Update Every 30 Minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=30s
Unit=world-cup-prediction-daily.service
Persistent=true

[Install]
WantedBy=timers.target
